# kubernetes/k8s-pod-recovery/files/k8s-pod-recovery.service
[Unit]
Description=Kubernetes Pods Recovery Service
After=kubelet.service
Requires=kubelet.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/k8s-pod-recovery start
ExecStop=/usr/local/sbin/k8s-pod-recovery stop
PIDFile=/var/run/k8s-pod-recovery.pid
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
