# devtools/nfscheck/recipes-common/nfscheck/files/nfscheck.service
[Unit]
Description=nfscheck
After=syslog.target network.target nfs-mountd.service sw-patch.service

[Service]
Type=simple
ExecStart=/bin/sh /usr/bin/nfscheck.sh

[Install]
WantedBy=multi-user.target
